# deploy/linux/definode-web.service
[Unit]
Description=DeFiNodeNexus Next.js Web Application
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=definode
Group=definode
WorkingDirectory=/opt/definode/DeFiNodeNexus
EnvironmentFile=/etc/definode/web.env
# standalone output: run server.js directly (no npm needed at runtime)
ExecStart=/usr/bin/node .next/standalone/server.js
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=definode-web
# Limit resources
LimitNOFILE=65536
Nice=0
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
