# deploy/linux/definode-health.timer
[Unit]
Description=Run DeFiNodeNexus daily health check

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true
Unit=definode-health.service

[Install]
WantedBy=timers.target