# deploy/package/linux/definode-health.service
[Unit]
Description=DeFiNodeNexus Daily Health Check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=definode
Group=definode
WorkingDirectory=/opt/definode/DeFiNodeNexus
EnvironmentFile=/etc/definode/keeper.env
ExecStart=/usr/bin/env node scripts/daily-health-check.js
TimeoutStartSec=300
Nice=10

[Install]
WantedBy=multi-user.target